# modbus2mqtt/rootfs/etc/s6-overlay/s6-rc.d/modbus2mqtt/run
#!/command/with-contenv bashio
set +e
cd /usr/app
export HASSIO_TOKEN
npx modbus2mqtt --yaml /data --ssl /ssl
