# modbus2mqtt.latest/docker/rootfs/etc/s6-overlay/s6-rc.d/modbus2mqtt/run
#!/command/with-contenv bashio
set +e
cd /usr/app
export HASSIO_TOKEN
inspect=""
if [ `bashio::addon.network | grep 9229 | wc -l` == 1 ]
then
  inspect=--inspect
fi
node $inspect dist/server/modbus2mqtt.js --yaml /data --ssl /ssl
 
